# scripts/launchd/com.kapshytar.adbfs-phone.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kapshytar.adbfs-phone</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__LAUNCHD_RUN_SH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Перезапускать только если процесс упал (exit != 0).
         Если телефон не найден — обёртка выходит 0, перезапуска нет. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/adbfs-phone.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/adbfs-phone.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__ADB_DIR__:/usr/local/bin:/usr/bin:/bin</string>
        <key>ADBFS_ROOT</key>
        <string>/storage/emulated/0</string>
    </dict>

</dict>
</plist>
